# ansible/roles/udiskie/files/hosts/common/.config/systemd/user/udiskie.service
[Unit]
Description=Udiskie
After=graphical-session.target


[Service]
Type=simple
ExecStart=/usr/bin/udiskie --appindicator -t
Restart=on-failure
Slice=background-graphical.slice

