# ansible/roles/swaysome/files/hosts/common/.config/systemd/user/swaysome.service
[Unit]
Description=Swaysome: AwesomeWM style workspaces for Sway
PartOf=graphical-session.target

[Service]
Type=oneshot
ExecStart=/usr/bin/swaysome init 1

