# ansible/roles/bluetooth/files/hosts/common/.config/systemd/user/blueman-applet.service
[Unit]
Description=Blueman Applet
After=graphical-session.target


[Service]
Type=simple
ExecStart=/usr/bin/blueman-applet
Restart=on-failure
Slice=background-graphical.slice

