# ansible/roles/hypridle/files/hosts/common/.config/systemd/user/hypridle.service
[Unit]
Description=Hypridle

[Service]
Type=exec
ExecStart=/usr/bin/hypridle
Slice=background-graphical.slice
Restart=on-failure
